# dycon_web/config/celery_beat.conf
[program:celerybeat]
directory=/home/landges/dycon/dycon_web
command=/home/landges/dycon/venv/bin/celery --app dycon_web beat --loglevel=INFO

user=landges
numprocs=1
stdout_logfile=/home/landges/dycon/dycon_web/logs/celery/beat-access.log
stderr_logfile=/home/landges/dycon/dycon_web/logs/celery/beat-error.log
stdout_logfile_maxbytes=50
stderr_logfile_maxbytes=50
stdout_logfile_backups=10
stderr_logfile_backups=10 
autostart=true
autorestart=true
startsecs=10

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999